# jenkins-master/configs/jenkins.ini
[program:jenkins]
command = JAVA_PATH JAVA_OPTS -jar /opt/jenkins.war --httpPort=8080 --ajp13Port=-1
stdout_logfile = syslog
redirect_stderr = true
user = APP_USER
