# was85/configs/was85.ini
[program:was85]
user=asrun
directory=/opt/was85/bin
command=/opt/run.sh
