# base/jira/configs/jira.ini
[program:jira]
command = /opt/jira/bin/start-jira.sh -fg
stdout_logfile = syslog
redirect_stderr = true
user = APP_USER
