# stash-base/configs/stash.ini
[program:stash]
command = /opt/stash/bin/start-stash.sh -fg
stdout_logfile = syslog
redirect_stderr = true
user = APP_USER
