# base/base7/jenkins-master/configs/jenkins.ini
[program:jenkins]
command = JAVA_PATH JAVA_OPTS -jar /opt/jenkins.war --httpPort=8050 --ajp13Port=9050
stdout_logfile = /tmp/jenkins_stdout.log
stderr_logfile = /tmp/jenkins_stderr.log
user = APP_USER
