# base/base8/sonar/configs/sonar.ini
[program:sonar]
command = /opt/sonar/bin/linux-x86-64/sonar.sh start
stdout_logfile = syslog
redirect_stderr = true
user = APP_USER
